# config/includes.chroot/etc/systemd/system/poweroff.target
[Unit]
Description=PH4NTXM Poweroff Panic
DefaultDependencies=no
Requires=ph4ntxm-panic.service
After=ph4ntxm-panic.service
AllowIsolate=yes
