# config/includes.chroot/etc/xdg/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network Manager
Exec=nm-applet --indicator
OnlyShowIn=XFCE;
X-GNOME-Autostart-enabled=true
