# config/includes.chroot/etc/xdg/autostart/ph4ntxm-panic-button-tray.desktop
[Desktop Entry]
Type=Application
Name=PH4NTXM Panic Button Tray
Exec=/usr/local/bin/ph4ntxm-panic-button-tray
Icon=ph4ntxm-panic-button
Comment=PH4NTXM Panic Button Tray
X-GNOME-Autostart-enabled=true
